# service/qthumbdump/run
#!/bin/sh

exec 2>&1

BIN=$(readlink -f -- $(dirname -- "$0"))

exec setuidgid xunxin.wan  $BIN/dump
